# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm7 -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:queue]
command=php /var/www/html/artisan queue:work --queue=high,default
startsecs=0
autorestart=false
startretries=1
logfile=/tmp/supervisord.log

[program:schedule]
command=/bin/sh -c 'while true; do date && php /var/www/html/artisan schedule:run; sleep 60; done'
stdout_logfile=/tmp/supervisord.log
autostart=true
autorestart=true
numprocs=1